# systemd/control-server-ttyd.service
[Unit]
Description=Control Server - Web Terminal (ttyd)
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/ttyd -p 7681 -i 0.0.0.0 /bin/bash
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
