# systemd/control-server-web-gui.service
[Unit]
Description=Control Server Web GUI - Web Management Panel
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/control-server-web-gui
Environment="PATH=/opt/control-server-web-gui/venv/bin:/usr/bin:/bin"
ExecStart=/opt/control-server-web-gui/venv/bin/python -m uvicorn app.main:app --host 0.0.0.0 --port 1206
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
