# systemd/websockify-vnc.service
[Unit]
Description=Websockify - VNC over WebSocket (port 6080)
After=network.target x11vnc.service
Requires=x11vnc.service

[Service]
Type=simple
# websockify có thể từ /usr/bin/websockify (apt) hoặc venv
ExecStart=/usr/bin/env websockify 0.0.0.0:6080 localhost:5900
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
